# supervisor/tornado.ini
[group:tornadoes]
programs=tornado-8000,tornado-10090

[program:tornado-8000]
command=python /apps/vin-decoder/h5/server.py
directory=/apps/vin-decoder/h5
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:tornado-10090]
command=python /apps/vin-decoder/app.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info
