# supervisor/tornado.ini
[program:tornado-8000]
command=python /apps/vin-decoder/h5/server.py
directory=/apps/vin-decoder/h5
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:tornado-10090]
command=python /apps/vin-decoder/app.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-31f]
command=python /apps/vin-decoder/spider/confirm.py 31f
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-kuaidaili]
command=python /apps/vin-decoder/spider/confirm.py kuaidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-xicidaili]
command=python /apps/vin-decoder/spider/confirm.py xicidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-66ip]
command=python /apps/vin-decoder/spider/confirm.py 66ip
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info
