# supervisor/vin-decoder.ini
[program:tornado-8000]
command=python /apps/vin-decoder/h5/server.py
directory=/apps/vin-decoder/h5
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:tornado-10090]
command=python /apps/vin-decoder/app.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-31f]
command=python /apps/vin-decoder/proxy/confirm.py 31f
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-kuaidaili]
command=python /apps/vin-decoder/proxy/confirm.py kuaidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-xicidaili]
command=python /apps/vin-decoder/proxy/confirm.py xicidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-66ip]
command=python /apps/vin-decoder/proxy/confirm.py 66ip
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-httpsdaili]
command=python /apps/vin-decoder/proxy/confirm.py httpsdaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-nianshao]
command=python /apps/vin-decoder/proxy/confirm.py nianshao
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info




[program:spider-proxy-31f]
command=python /apps/vin-decoder/spider/proxy/cn_31f.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-66ip]
command=python /apps/vin-decoder/spider/proxy/cn_66ip.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-kuaidaili]
command=python /apps/vin-decoder/spider/proxy/com_kuaidaili.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-xicidaili]
command=python /apps/vin-decoder/spider/proxy/com_xicidaili.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-httpsdaili]
command=python /apps/vin-decoder/spider/proxy/com_httpsdaili.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info


[program:anonymous-spider-proxy-31f]
command=python /apps/vin-decoder/spider/proxy/anonymous/31f.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:anonymous-spider-proxy-66ip]
command=python /apps/vin-decoder/spider/proxy/anonymous/66ip.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:anonymous-spider-proxy-kuaidaili]
command=python /apps/vin-decoder/spider/proxy/anonymous/kuaidaili.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:anonymous-spider-proxy-xicidaili]
command=python /apps/vin-decoder/spider/proxy/anonymous/xicidaili.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:anonymous-spider-proxy-httpsdaili]
command=python /apps/vin-decoder/spider/proxy/anonymous/httpsdaili.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:anonymous-spider-proxy-nianshao]
command=python /apps/vin-decoder/spider/proxy/anonymous/nianshao.py
directory=/apps/vin-decoder
user=root
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info
