# supervisor/vin.ini
[program:tornado-8000]
command=python /apps/vin-decoder/h5/server.py
directory=/apps/vin-decoder/h5
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:tornado-10090]
command=python /apps/vin-decoder/app.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-proxy-31f]
command=python /apps/vin-decoder/proxy/confirm.py 31f
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-proxy-kuaidaili]
command=python /apps/vin-decoder/proxy/confirm.py kuaidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-proxy-xicidaili]
command=python /apps/vin-decoder/proxy/confirm.py xicidaili
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:confirm-proxy-66ip]
command=python /apps/vin-decoder/proxy/confirm.py 66ip
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:spider-proxy-31f]
command=python /apps/vin-decoder/spider/proxy/31f.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:spider-proxy-66ip]
command=python /apps/vin-decoder/spider/proxy/66ip.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:spider-proxy-kuaidaili]
command=python /apps/vin-decoder/spider/proxy/kuaidaili.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info

[program:spider-proxy-xicidaili]
command=python /apps/vin-decoder/spider/proxy/xicidaili.py
directory=/apps/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=info
